# work/free_bird/rpmtop/BUILD/freebird-1.0.0/freebird.service
[Unit]
Description=Freed Bird Service

[Service]
Type=simple
ExecStart=/usr/bin/free_bird.py
