# work/trunk/systemd/puma-watchdog.service
[Unit]
Description=Puma Service WatchDog Service

[Service]
Type=simple
ExecStart=/opt/puma/puma-watchdog
