# work/free_bird/freebird.timer
[Unit]
Description= Freedom Bird Timer

[Timer]
OnBootSec=5min
OnUnitActiveSec=5min
Unit=freedom_bird.service

[Install]
WantedBy=multi-user.target
