# work/free_bird/freebird.service
[Unit]
Description=Freed Bird Service

[Service]
Type=forking
ExecStart=/usr/bin/freebird.py
